# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mcp_server]
command=uv run python src/mcp_server.py
directory=/app
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/mcp_server.log
stderr_logfile=/var/log/supervisor/mcp_server.err.log

[program:voice_server]
command=uv run python src/voice_server.py dev
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/var/log/supervisor/voice_server.log
stderr_logfile=/var/log/supervisor/voice_server.err.log

[program:frontend_server]
command=uv run python src/frontend_server.py
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/frontend_server.log
stderr_logfile=/var/log/supervisor/frontend_server.err.log
